# rc.local
#!/bin/bash
#
# /etc/rc.local: multi-user startup script
#

# end of file
